# embedtools/rtcsync.service
[Unit]
Description= Writes sys time to RTC, if we are currently synced

[Service]
Type=oneshot
ExecStart=/usr/bin/rtcsync.sh